# crontab
# SEM Assistant Elisp Daemon Crontab
# Install with: COPY crontab /etc/cron.d/sem-cron
# Timezone environment lines are injected at container startup by dev/start-cron.

# Inbox processing - every 30 minutes
*/30 * * * *	root	emacsclient -s sem-server -e "(sem-core-process-inbox)"

# Git pre-pull reconciliation - every 5 minutes
*/5 * * * *	root	emacsclient -s sem-server -e "(sem-git-sync-org-roam-prepull)"

# Inbox purge - daily at 4:00 AM (only time inbox-mobile.org is written)
0 4 * * *	root	emacsclient -s sem-server -e "(sem-core-purge-inbox)"

# Elfeed update - four times before digest generation (5AM, 6AM, 7AM, 8AM)
0 5 * * *	root	emacsclient -s sem-server -e "(sem-rss-elfeed-update-cron)"
0 6 * * *	root	emacsclient -s sem-server -e "(sem-rss-elfeed-update-cron)"
0 7 * * *	root	emacsclient -s sem-server -e "(sem-rss-elfeed-update-cron)"
0 8 * * *	root	emacsclient -s sem-server -e "(sem-rss-elfeed-update-cron)"

# RSS morning digest - daily at 9:30 AM (after all elfeed updates complete)
30 9 * * *	root	emacsclient -s sem-server -e "(sem-rss-generate-morning-digest)"

# Daemon readiness watchdog - every 15 minutes
*/15 * * * *	root	/usr/local/bin/sem-daemon-watchdog

# GitHub sync for org-roam - every 3 hours (00:00, 03:00, 06:00, 09:00, 12:00, 15:00, 18:00, 21:00)
0 */3 * * *	root	emacsclient -s sem-server -e "(sem-git-sync-org-roam)"
